# scripts/rice-ota.service
[Unit]
Description=Pull latest Hyprland rice from GitHub

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/ota-update.sh
